# docker/containers/supervisor/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisor/supervisord.pid
logfile=/dev/stdout
logfile_maxbytes=0

[program:queue]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --timeout=90
user=www-data
autostart=true
autorestart=true
priority=20
numprocs=1
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:schedule]
command=php /var/www/html/artisan schedule:work
user=www-data
autostart=true
autorestart=true
priority=20
numprocs=1
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
